# neovim/rc.local
set relativenumber

" Fix guioptions for Neovim (not supported)
if has('gui_running')
  set guioptions=egmrti
endif

" Normal mode in terminal.
tnoremap <Esc> <C-\><C-n>

" Enable mouse support in terminal
set mouse=a

" Always show cursor postion
set ruler

" Hide the buffer when abandoned
set hidden

" Keep some lines visible when scrolling
set scrolloff=3

" Maximum width of text that is being inserted
set textwidth=79

" Highlight past 80 columns
set colorcolumn=80
"
" Automatically load a file changed outside
set autoread

" Enable autoindent
set autoindent

" Fast edit .vimrc file
nmap <localleader>e :e! $MYVIMRC<cr>

" Change backup directory
set backupdir=~/.tmp

" Change the directory of swp files
set directory=~/.tmp

" Faster and improved grep using ag (The Silver Searcher)
if executable('ag')
    set grepprg=ag\ --vimgrep\ $*
    se\ grepformat=%f:%l:%c:%m
endif

let g:winresizer_vert_resize=1
let g:winresizer_horiz_resize=1

let g:easyescape_chars = { "j": 2 }
let g:easyescape_timeout = 500

cnoremap jj <ESC>
 
" enable AutoSave on Vim startup
let g:auto_save = 1
let g:auto_save_events = ["InsertLeave", "TextChanged"]

" Vim airline
let g:airline#extensions#tabline#formatter = 'unique_tail_improved' 
let g:airline_theme='molokai'
let g:airline_powerline_fonts = 1

let g:go_rename_command = 'gopls'

let g:tagbar_type_go = {
	\ 'ctagstype' : 'go',
	\ 'kinds'     : [
		\ 'p:package',
		\ 'i:imports:1',
		\ 'c:constants',
		\ 'v:variables',
		\ 't:types',
		\ 'n:interfaces',
		\ 'w:fields',
		\ 'e:embedded',
		\ 'm:methods',
		\ 'r:constructor',
		\ 'f:functions'
	\ ],
	\ 'sro' : '.',
	\ 'kind2scope' : {
		\ 't' : 'ctype',
		\ 'n' : 'ntype'
	\ },
	\ 'scope2kind' : {
		\ 'ctype' : 't',
		\ 'ntype' : 'n'
	\ },
	\ 'ctagsbin'  : 'gotags',
	\ 'ctagsargs' : '-sort -silent'
\ }

set autowrite

map <C-n> :cnext<CR>
map <C-m> :cprevious<CR>
nnoremap <leader>a :cclose<CR>

let g:go_list_type = "quickfix"
autocmd FileType go nmap <leader>t  <Plug>(go-test)
let g:go_fmt_command = "goimports"
let g:go_highlight_types = 1
let g:go_highlight_fields = 1
let g:go_highlight_functions = 1
let g:go_highlight_function_calls = 1
let g:go_highlight_build_constraints = 1

let g:rehash256 = 1
let g:molokai_original = 1

let g:coc_global_extensions = ['coc-tsserver']

if isdirectory('./node_modules') && isdirectory('./node_modules/prettier')
  let g:coc_global_extensions += ['coc-prettier']
endif

if isdirectory('./node_modules') && isdirectory('./node_modules/eslint')
  let g:coc_global_extensions += ['coc-eslint']
endif

nmap <leader>rn <Plug>(coc-rename)

let g:nerdtree_tabs_open_on_gui_startup = 0

"" Close buffer
noremap <leader>c :bd!<CR>

" Terminal Function
let g:term_buf = 0
let g:term_win = 0
function! TermToggle(height)
    if win_gotoid(g:term_win)
        hide
    else
        botright new
        exec "resize " . a:height
        try
            exec "buffer " . g:term_buf
        catch
            call termopen($SHELL, {"detach": 0})
            let g:term_buf = bufnr("")
            set nonumber
            set norelativenumber
            set signcolumn=no
        endtry
        startinsert!
        let g:term_win = win_getid()
    endif
endfunction

" Terminal go back to normal mode
tnoremap <Esc> <C-\><C-n>
tnoremap :q! <C-\><C-n>:q!<CR>

" terminal emulation
nnoremap <silent> <leader>sh :call TermToggle(10)<CR>

let g:deoplete#enable_at_startup = 1

if exists('*deoplete#custom#option')
  call deoplete#custom#option('omni_patterns', { 'go': '[^. *\t]\.\w*' })
endif

let g:lastplace_ignore = "gitcommit,gitrebase,svn,hgcommit"

let g:lastplace_ignore_buftype = "quickfix,nofile,help"

let NERDTreeShowHidden=1

let g:go_metalinter_command='golangci-lint'
let g:go_metalinter_autosave = 0
let g:go_metalinter_enabled = [ 
    \"asasalint",
    \"asciicheck",
    \"bidichk",
    \"bodyclose",
    \"containedctx",
    \"contextcheck",
    \"cyclop",
    \"decorder",
    \"depguard",
    \"dogsled",
    \"dupl",
    \"dupword",
    \"durationcheck",
    \"errcheck",
    \"errname",
    \"errorlint",
    \"execinquery",
    \"exhaustive",
    \"exportloopref",
    \"forbidigo",
    \"forcetypeassert",
    \"gci",
    \"gochecknoinits",
    \"gocognit",
    \"goconst",
    \"gocritic",
    \"gocyclo",
    \"godot",
    \"godox",
    \"goerr113",
    \"gofmt",
    \"goheader",
    \"goimports",
    \"gomnd",
    \"gomoddirectives",
    \"gomodguard",
    \"goprintffuncname",
    \"gosimple",
    \"govet",
    \"grouper",
    \"importas",
    \"ineffassign",
    \"interfacebloat",
    \"lll",
    \"loggercheck",
    \"maintidx",
    \"makezero",
    \"nakedret",
    \"nestif",
    \"nilerr",
    \"nilnil",
    \"nlreturn",
    \"noctx",
    \"nolintlint",
    \"nonamedreturns",
    \"nosprintfhostport",
    \"paralleltest",
    \"prealloc",
    \"predeclared",
    \"promlinter",
    \"reassign",
    \"revive",
    \"staticcheck",
    \"stylecheck",
    \"tenv",
    \"testableexamples",
    \"testpackage",
    \"thelper",
    \"tparallel",
    \"typecheck",
    \"unconvert",
    \"unparam",
    \"unused",
    \"usestdlibvars",
    \"varnamelen",
    \"whitespace",
    \"wrapcheck",
    \"wsl",
    \]


set spelllang=en,pt_br
set spellsuggest=best,9

" Auto-download spell files if missing
if !filereadable(expand('~/.config/nvim/spell/pt.utf-8.spl'))
  autocmd VimEnter * silent! :mkspell! ~/.config/nvim/spell/pt.utf-8.add
endif

" let g:go_debug = ['shell-commands']
let g:auto_save = 1  " enable AutoSave on Vim startup

" Find files using Telescope command-line sugar.
nnoremap <leader>ff <cmd>Telescope find_files<cr>
nnoremap <leader>fg <cmd>Telescope live_grep<cr>
nnoremap <leader>fb <cmd>Telescope buffers<cr>
nnoremap <leader>fh <cmd>Telescope help_tags<cr>


" ============================================
" Writing & Markdown Configuration
" ============================================

" Spell check configuration
set spell
set spelllang=en,pt_br
set spellsuggest=best,9

" Goyo (distraction-free writing)
nnoremap <leader>gy :Goyo<CR>

" Limelight (focus on current paragraph)
nnoremap <leader>ll :Limelight!!<CR>

" Goyo integration with Limelight
autocmd! User GoyoEnter Limelight
autocmd! User GoyoLeave Limelight!

" Limelight configuration
let g:limelight_conceal_ctermfg = 'gray'
let g:limelight_conceal_guifg = '#777777'
let g:limelight_default_coefficient = 0.7

" Vim-pencil configuration for writing
let g:pencil#wrapModeDefault = 'soft'
augroup pencil
  autocmd!
  autocmd FileType markdown,mkd call pencil#init()
  autocmd FileType text         call pencil#init()
augroup END

" Markdown configuration
let g:vim_markdown_folding_disabled = 1
let g:vim_markdown_conceal = 0
let g:vim_markdown_frontmatter = 1
let g:vim_markdown_toml_frontmatter = 1
let g:vim_markdown_json_frontmatter = 1

" Markdown shortcuts
autocmd FileType markdown nnoremap <buffer> <leader>p :MarkdownPreview<CR>
autocmd FileType markdown setlocal spell
autocmd FileType markdown setlocal linebreak

" Quick spell fix
nnoremap <leader>z 1z=

" Toggle spell check
nnoremap <leader>sp :setlocal spell!<CR>
